# linux/adhoc-backups/backup_adhoc_remotes.cron
#!/usr/bin/env bash
# vim: shiftwidth=2 tabstop=2 expandtab
set -euo pipefail
[[ "${TRACE-0}" == "1" ]] && set -x

# can be installed as /etc/cron.hourly/backup_adhoc_remotes.cron

name=$(basename -- "$BASH_SOURCE")
lockfile=/var/lock/$name

(
  flock -x -w 10 200 || exit 0

  shopt -s nullglob
  for backup in /Volumes/backup-*/backup--*/backup.marker
  do
    "$backup"
  done > /var/tmp/$name.$(date -u +%Y-%m-%d.%H%M%Sz).log

) 200>"$lockfile"
